# guardian.service
[Unit]
Description=Guardian EDR Agent
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/guardian
ExecStart=/usr/bin/python3 /home/%i/guardian/guardian.py --headless
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
